# qr3_93.conf
[program:qr3_93]
command=python3 /home/curator3/qr3_93/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_93/test.err.log
stdout_logfile=/home/curator3/qr3_93/test.err.log
